# protocolling/.config/systemd/user/protocolling_sqlite.timer
[Unit]
Description=Run protocolling hourly

[Timer]
OnCalendar=*-*-* *:00,30:00
Persistent=True

[Install]
WantedBy=default.target
